# watchdog/watchdog.service
[Unit]
Description=Watchdog

[Service]
ExecStart=/opt/watchdog/watchdog /dev/watchdog0

[Install]
WantedBy=multi-user.target
