# watchdog/watchdog.service
[Unit]
Description=Watchdog

[Service]
ExecStart=python3 /opt/watchdog.py

[Install]
WantedBy=multi-user.target
